# crontab
*/15 * * * * python3 /home/pi/mitemp/mitemp-mqtt.py 58:2D:34:xx:xx:xx -m -e -s mqtt-server