# crontab
*/15 * * * * python3 /home/pi/apps/mitemp/mitemp_mqtt.py 58:2D:34:36:xx:xx -m -e -s mqtt-server